# Misc/onlogin.service
[Unit]
Description=Run startup script

[Service]
ExecStart=/home/stuart/Dotfiles/startup.sh

[Install]
WantedBy=default.target
