# Misc/onlogin.service
[Unit]
Description=Run startup script

[Service]
ExecStart=$HOME/.startup.sh

[Install]
WantedBy=multi-user.target
